# config/autostart/i3.desktop
# pets: package=xcfe4
# pets: package=i3
# pets: symlink=~/.gitconfig
#
# Warning! This file has been generated by pets(1). Any manual modification
# will be lost.

[Desktop Entry]
Version=0.9.4
Type=Application
Name=i3
Comment=The only way to use a computer
Exec=i3
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false
